# request_handler.service
[Unit]
Description=HTTP Request Handling Service

[Service]
ExecStart=/usr/bin/python3 /home/sanya/systemd_echo/request_handler.py